# overlay/etc/udev/rules.d/z99-innobox.rules
# Dump cool stuff on the usb stick
KERNEL=="sd*", BUS=="usb", ACTION=="add", RUN+="/usr/share/innobox-dump/innobox-dump.py %k"
